# infra/systemd/telegram-codex-shutdown.service
[Unit]
Description=Power off the Telegram Codex worker VM after a drained run

[Service]
Type=oneshot
ExecStart=/sbin/shutdown -h now
