# infra/systemd/telegram-codex-shutdown.path
[Unit]
Description=Watch for a Telegram Codex worker shutdown request

[Path]
PathExists=/run/telegram-codex-worker/shutdown-request
Unit=telegram-codex-shutdown.service

[Install]
WantedBy=multi-user.target
